# .config/systemd/user/udiskie.service
[Unit]
Description=Start udiskie

[Service]
Type=simple
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
