# .config/systemd/user/warn-battery.service
[Unit]
Description=Notify on low battery

[Service]
Type=oneshot
ExecStart=/usr/local/bin/warn-battery

[Install]
WantedBy=default.target
