# .crontab
# wallpaper
*/60 * * * * curl -sS 'https://unsplash.it/3840/2160/?random' | DISPLAY=:0.0 feh --bg-scale -

# low battery popup
*/1 * * * * DBUS_SESSION_BUS_ADDRESS=$(grep -z DBUS_SESSION_BUS_ADDRESS /proc/$(pidof i3)/environ | cut -d= -f2-) /usr/local/bin/warn-battery

# installed packages
0 * * * * yaourt -Qe --nocolor | cut -d' ' -f1,3- | sort > $HOME/.installed-packages-arch-$(hostname)
0 * * * * python $HOME/.local/bin/pipdeptree --warn silence | grep -E '^\w+' | cut -d'=' -f1 | sort> $HOME/.installed-packages-python-$(hostname)
0 * * * * python2 $HOME/.local/bin/pipdeptree --warn silence | grep -E '^\w+' | cut -d'=' -f1 | sort > $HOME/.installed-packages-python2-$(hostname)
0 * * * * npm list -g -depth 0 | grep -v /usr | sed '/^$/d' | cut -d' ' -f2 | cut -d'@' -f1 | sort > $HOME/.installed-packages-node-$(hostname)
