# .crontab
0 * * * * DISPLAY=:0.0 feh --bg-scale "$(find ~/Documents/interfacelift-downloader/wallpapers/*.jpg | shuf -n1)"
*/1 * * * * DBUS_SESSION_BUS_ADDRESS=$(grep -z DBUS_SESSION_BUS_ADDRESS /proc/$(pidof i3)/environ | cut -d= -f2-) /usr/local/bin/warn-battery

# installed packages
0 * * * * yaourt -Qe --nocolor | cut -d' ' -f1,3- | sort > $HOME/.installed-packages-arch-$(hostname)
0 * * * * python $(which pipdeptree) --warn silence | grep -E '^\w+' | cut -d'=' -f1 | sort> $HOME/.installed-packages-python-$(hostname)
0 * * * * python2 $(which pipdeptree) --warn silence | grep -E '^\w+' | cut -d'=' -f1 | sort > $HOME/.installed-packages-python2-$(hostname)
0 * * * * npm list -g -depth 0 | grep -v /usr | sed '/^$/d' | cut -d' ' -f2 | cut -d'@' -f1 | sort > $HOME/.installed-packages-node-$(hostname)
*/60 * * * * curl -sS 'https://unsplash.it/1920/1080/?random' > $HOME/.wallpaper.jpg && DISPLAY=:0.0 feh --bg-scale $HOME/.wallpaper.jpg
