# .crontab
0 * * * * DISPLAY=:0.0 feh --bg-scale "$(find ~/Documents/interfacelift-downloader/wallpapers/*.jpg | shuf -n1)"
*/5 * * * * warn-battery
0 * * * * sudo reflector --latest 200 --protocol http --sort rate --save /etc/pacman.d/mirrorlist
0 * * * * yaourt -Qe --nocolor | cut -d ' ' -f 1,3- | sort > $HOME/.installed-packages
