# etc/rc.local
# Fix touch keys
/usr/lib/udev/keymap input/event0 /opt/dell_touchkeys_keymap
