# etc/udev/rules.d/mouse.rules
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="hidpp_battery_0", RUN+="/usr/bin/solaar config 1 dpi 600"
